# etc/profile.d/lah.sh
alias cp='/bin/cp -i'
alias ll='/bin/ls -alFh'
alias rm='/bin/rm -i'
